# src/mirror/nix-mirror@nixos-19.03-small.timer
[Unit]
Description=Mirror nixos-19.03-small once every eight hours

[Timer]
OnCalendar=0/8:00:00
Persistent=true

[Install]
WantedBy=timers.target
